# bootstraps_nyz/bootstrap-aarch64/usr/var/service/telnetd/run
#!/data/data/com.bingduoduo/files/usr/bin/sh
exec telnetd -F
